# root/etc/services.d/wikijs/run
#!/usr/bin/with-contenv bash

cd /app/wiki || exit

exec \
	s6-setuidgid abc /usr/bin/node server
